# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dashboard-backend]
command=node /app/dashboard/backend/server.js
directory=/app/dashboard/backend
autostart=true
autorestart=true
stderr_logfile=/app/logs/dashboard-backend.err.log
stdout_logfile=/app/logs/dashboard-backend.out.log
environment=NODE_ENV="production",PORT="3456"

[program:luna-server]
command=node /app/dashboard/backend/luna-server.js
directory=/app/dashboard/backend
autostart=true
autorestart=true
stderr_logfile=/app/logs/luna-server.err.log
stdout_logfile=/app/logs/luna-server.out.log
environment=NODE_ENV="production",LUNA_PORT="3458"
